# deploy/systemd/tdi_reader.service
[Unit]
Description=FTDI to CRSF bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tdi_reader
Group=dialout
WorkingDirectory=__PROJECT_ROOT__
ExecStart=__PROJECT_ROOT__/build-linux/tdi_reader __CONFIG_PATH__
Restart=always
RestartSec=2
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
